# deploy/supervisor/smkmuhbligo-worker.conf
; =============================================================================
; Supervisor Configuration — SMK Muhammadiyah Bligo
; Path  : /etc/supervisor/conf.d/smkmuhbligo.conf
; App   : /var/www/smkmuhbligo
; Queue : Redis (ai, default)
; Updated: 2026-03-06
; =============================================================================

[program:smkmuhbligo-worker]

; ─── Command ─────────────────────────────────────────────────────────────────
; --queue=ai,default  : prioritaskan queue AI, kemudian default
; --tries=3           : max retry jika job gagal sebelum dianggap failed
; --max-time=3600     : restart worker otomatis setelah 1 jam (cegah memory leak)
; --sleep=3           : idle 3 detik jika queue kosong
command=php /var/www/smkmuhbligo/artisan queue:work redis --queue=ai,default --tries=3 --max-time=3600 --sleep=3

; ─── Working Directory ───────────────────────────────────────────────────────
directory=/var/www/smkmuhbligo

; ─── Process Identity ────────────────────────────────────────────────────────
user=www-data

; ─── Multiple Processes ──────────────────────────────────────────────────────
; 2 worker: prioritas AI dijamin melalui --queue=ai,default di kedua proses
numprocs=2
process_name=%(program_name)s_%(process_num)02d

; ─── Autostart & Autorestart ──────────────────────────────────────────────────
autostart=true
autorestart=unexpected
exitcodes=0,2

; ─── Graceful Shutdown ────────────────────────────────────────────────────────
; Tunggu job selesai sebelum kill (sesuaikan dengan --max-time)
stopasgroup=true
killasgroup=true
stopwaitsecs=3600
stopsignal=TERM

; ─── Logging ─────────────────────────────────────────────────────────────────
redirect_stderr=true
stdout_logfile=/var/log/smkmuhbligo-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/smkmuhbligo-worker-error.log
stderr_logfile_maxbytes=10MB
